# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:api]
command=uvicorn src.serving.api:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.err
environment=PYTHONPATH="/app",LOG_LEVEL="INFO",LOG_FORMAT="json"
startretries=3
startsecs=5

[program:ui]
command=streamlit run src/ui/app.py --server.port 7860 --server.address 0.0.0.0 --server.headless true --server.enableCORS false --server.enableXsrfProtection false
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ui.log
stderr_logfile=/var/log/supervisor/ui.err
environment=PYTHONPATH="/app",API_URL="http://localhost:8000"
startretries=3
startsecs=10
